# redis/redis-cluster/supervisor_redis_cluster.ini
[program:redis_cluster]
command=redis-server redis.conf
process_name=%(program_name)s_%(process_num)02d
directory=/home/wangx/github/linux-reference/redis/redis-cluster/cluster-test/70%(process_num)02d
numprocs=3
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/home/wangx/venv/bin/:/home/wangx/.local/bin:/home/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/home/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/home/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log
