# redis/redis-cluster/supervisor_redis_cluster_mac.ini
[program:redis_cluster]
process_name=redis_cluster_700%(process_num)s
numprocs=8
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/700%(process_num)s
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/700%(process_num)s/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/700%(process_num)s/stderr.log
