# redis/redis-cluster/supervisor_redis_cluster_mac.ini
[program:redis_cluster_7000]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7000
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7001]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7001
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7002]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7002
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7003]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7003
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7004]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7004
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7005]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7005
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7006]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7006
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7007]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7007
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

[program:redis_cluster_7008]
command=redis-server redis.conf
directory=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/7008
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/opt/homebrew/bin:/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stdout.log
stderr_logfile=/Users/wangx/github/linux-reference/redis/redis-cluster/cluster-test/stderr.log

